# Cave game/cave-6736772.20190622T103148719Z.full/data/etc/profile.d/mysql.sh
export MYSQL_USERNAME="letonus"
export MYSQL_PASSWORD="39kHjc3RQFthijJO"
export MYSQL_HOSTNAME="localhost"
